# crontab
*/30 * * * * python /app/main.py >> /dev/stdout
0 */6 * * * python /app/auth.py >> /dev/stdout