# crontab
*/10 * * * * python /app/main.py >> /dev/stdout
0 */12 * * * python /app/auth.py >> /dev/stdout